# config/project360.conf
[program:project360]
command=/home/john/venv/bin/gunicorn three_sixty_project.wsgi:application -c /home/john/three_sixty_project/config/gunicorn.conf.py
directory=/home/john/three_sixty_project
user=john
autorestart=true
redirect_stderr=true
stdout_logfile = /home/john/three_sixty_project/logs/debug.log