# config/project360.conf
[program:project360]
command=/home/www/code/project360/env/bin/gunicorn three_sixty_project.wsgi:application -c /home/www/code/project360/three_sixty_project/config/gunicorn.conf.py
directory=/home/www/code/project360/three_sixty_project
user=www
autorestart=true
redirect_stderr=true
